# projects/dns.palladius.eu/files/etc/cron.d/pull-dns-and-reload-if-good
*/30 * * * *    root     cd /etc/bind && git pull && service bind9 reload && touch bind-reloaded-last.touch
